# rpi/frontend/patobot_frontend.conf
[program:patobot_frontend]
autorestart=true
command=/home/pi/patobot/frontend/env/bin/python -m SimpleHTTPServer
directory=/home/pi/patobot/frontend/
user:pi
