# rpi/backend/patobot_backend.conf
[program:patobot_backend]
autorestart=true
command=/home/pi/patobot/backend/env/bin/gunicorn -b 0.0.0.0:5000 patobot:app
directory=/home/pi/patobot/backend/
user:pi
