# rpi/backend/patobot.conf
[program:patobot]
autorestart=true
command=/home/pi/patobot/backend/env/bin/gunicorn patobot:app
directory=/home/pi/patobot/backend/
user:pi
